# webstack/load_balancer/load_balancer.service
[Unit]
Description=Load Balancer

[Service]
ExecStart=/deploy/load_balancer/load_balancer
WorkingDirectory=/deploy/load_balancer
Restart=always
